# services/rgpio_monitor/run
#!/bin/sh
echo "*** starting RemoteGPIO ***"
exec 2>&1
exec nice -n 19 /data/RemoteGPIO/rgpio_monitor
