# services/rgpio_driver/run
#!/bin/sh
echo "*** starting RemoteGPIO ***"
exec 2>&1
exec /data/RemoteGPIO/rgpio_driver
